# scripts/csff/csff-trade-scanner.timer
[Unit]
Description=CSFF trade scanner — every 10 min during US market hours (15:35-21:55 Paris)

[Timer]
OnCalendar=Mon..Fri *:0/10:00 Europe/Paris
Persistent=false

[Install]
WantedBy=timers.target
